# layers/meta-vtu/recipes-vtu/vtu-telemetry/files/vtu-telemetry.service
[Unit]
Description=VTU MQTT Telemetry Publisher
Documentation=https://github.com/almirmujanovic/vtu-project
After=network.target vtu-ecu-sim.service
Wants=vtu-ecu-sim.service

[Service]
Type=simple
# Default broker - override with drop-in file for production
ExecStart=/usr/bin/vtu-telemetry -b tcp://localhost:1883 -i vcan0
Restart=on-failure
RestartSec=10

# Security hardening
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target