# layers/meta-vtu/recipes-vtu/vtu-obdgw/files/vtu-obdgw.service
[Unit]
Description=VTU OBD-II Gateway
Documentation=https://github.com/almirmujanovic/vtu-project
After=network.target vtu-ecu-sim.service
Wants=vtu-ecu-sim.service

[Service]
Type=simple
ExecStart=/usr/bin/vtu-obdgw vcan0
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target